# docs/examples/nginx_gunicorn/supervisor.conf
[program:nsupdate.info]
command = /srv/nsupdate/bin/gunicorn_start  ; Command to start app
user = www-data                             ; User to run as
stdout_logfile = /srv/nsupdate/logs/gunicorn_supervisor.log  ; Where to write log messages
redirect_stderr = true                      ; Save stderr in the same log
